# deploy/systemd/hermes-redis-gateway-worker.service
[Unit]
Description=Hermes Redis Gateway Worker
After=network-online.target redis-server.service
Wants=network-online.target

[Service]
Type=simple
User=song
Group=song
WorkingDirectory=/home/song/Desktop/hermes-redis-gateway
EnvironmentFile=/home/song/Desktop/hermes-redis-gateway/.env
ExecStart=/home/song/Desktop/hermes-redis-gateway/.venv/bin/hermes-redis-gateway-worker
Restart=always
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
